# _scripts/db-visualizer.service
[Unit]
Description=DB Visualizer Application
After=network.target

[Service]
Type=simple
User=kim
WorkingDirectory=/home/kim/CG3_Tech/Projects/Personal Project 2025/db-visualizer
ExecStart=/usr/local/bin/db-visualizer-start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (OPTIONAL - only needed if you want to override .env values)
# Your backend .env file will be automatically loaded by NestJS ConfigModule
# Uncomment these only if you need to override the .env file values:
# Environment="NODE_ENV=production"
# Environment="PORT=3000"

[Install]
WantedBy=multi-user.target

